# services/crosscompute-homepage.timer
[Timer]
OnCalendar=*-*-* *:00:00
Unit=crosscompute-homepage.service
