# services/crosscompute-homepage.service
[Service]
Type=oneshot
ExecStart=/home/crosscompute/Projects/crosscompute-homepage/scripts/build.sh
