# services/crosscompute-homepage.service
[Service]
Type=oneshot
ExecStart=/root/Projects/crosscompute-homepage/scripts/build.sh
